# deploy/sa-home-bot.service
# Пользовательский юнит: ~/.config/systemd/user/sa-home-bot.service
# systemctl --user enable --now sa-home-bot
[Unit]
Description=sa-home-bot (Telegram-фронтенд мониторинга)
# network-online.target в user-режиме недоступен, поэтому готовность сети ждём
# через ExecStartPre (ниже): бот при старте проверяет подписки и, если Telegram
# недоступен (DNS ещё не поднялся после ребута), помечает их broken до перезапуска.
After=network.target sa-home-monitor.service
# Бот заказывает монитор, но живёт и без него (переподключается сам).
Wants=sa-home-monitor.service

[Service]
Type=simple
WorkingDirectory=/home/sevboa/Github/sa-home-bot
# Гейт: не запускать бота, пока не резолвится api.telegram.org (значит сеть+DNS
# поднялись). Цикл ограничен TimeoutStartSec; при таймауте старт падает и Restart
# пробует снова. Это устраняет «broken-подписки после отключения электричества».
ExecStartPre=/bin/sh -c 'until getent hosts api.telegram.org >/dev/null 2>&1; do sleep 3; done'
ExecStart=/home/sevboa/Github/sa-home-bot/.venv/bin/sa-home-bot --config /home/sevboa/Github/sa-home-bot/config.toml
Restart=on-failure
RestartSec=5
TimeoutStartSec=300
KillSignal=SIGTERM
TimeoutStopSec=30
Environment=PATH=/home/sevboa/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

[Install]
WantedBy=default.target
